# resources/deb/variants/debian/DEBIAN/postinst
#!/bin/sh

systemctl enable depot.service