# deb/lib/systemd/system/depot.service
[Unit]
Description=Depot

[Service]
Type=simple
ExecStart=/etc/init.d/havis-depot run
Restart=on-failure

[Install]
WantedBy=multi-user.target