# resources/deb/variants/debian/lib/systemd/system/depot.service
[Unit]
Description=Depot

[Service]
Type=simple
ExecStart=/etc/init.d/httpd run
Restart=on-failure

[Install]
WantedBy=multi-user.target