# debian/postinst
#!/bin/sh

if [ -d /run/systemd/system ]; then
  systemctl enable depot.service || true
fi